# _next/data/tTH-BFv6YqOvxwhFqYTb1/rules/sigma-6c4e2f43-d94d-4ead-b64d-97e53fa2bd05.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2021-10-15T00:00:00.000Z","creation_date":"2021-10-15T00:00:00.000Z","source_siem":"Sigma","source_siem_id":"sigma","source_url":"https://github.com/SigmaHQ/sigma/blob/master/rules/linux/file_event/file_event_lnx_susp_cron_file_created.yml"},"rule":{"name":"New Cron File Created","description":"Detects the creation of cron files in Cron directories, which could indicate potential persistence mechanisms being established by an attacker.\nNote that not all cron file creations are malicious - legitimate system administration activities and software installations may also create cron files.\nThis detection should be investigated in context, considering factors such as the user creating the file, the timing of creation, and the contents of the cron job.\nFocus investigation on unexpected cron files created by non-administrative users or during suspicious timeframes.\nAdditionally, it is recommended to review the contents of the newly created cron files to assess their intent.\nFurthermore, it is suggested to baseline normal cron file creation and apply additional filters to reduce false positives based on the specific environment.\n","tags":["Rule Type: Sigma","Data Source: Cross-SIEM","Tactic: Execution","Tactic: Persistence","Technique: T1053.003","OS: Linux","Domain: file event","SIEM: Sigma"],"severity":"low","risk_score":21,"type":"sigma","language":"sigma","query":"{\n  \"selection_cron_dirs\": {\n    \"TargetFilename|startswith\": [\n      \"/etc/cron.d/\",\n      \"/etc/cron.daily/\",\n      \"/etc/cron.hourly/\",\n      \"/etc/cron.monthly/\",\n      \"/etc/cron.weekly/\",\n      \"/var/spool/cron/crontabs/\",\n      \"/var/spool/cron/root\"\n    ]\n  },\n  \"selection_cron_special_files\": {\n    \"TargetFilename|contains\": [\n      \"/etc/cron.allow\",\n      \"/etc/cron.deny\",\n      \"/etc/crontab\"\n    ]\n  },\n  \"filter_optional_legit_cron\": {\n    \"TargetFilename\": [\n      \"/etc/cron.daily/apt\",\n      \"/etc/cron.daily/dpkg\",\n      \"/etc/cron.daily/passwd\",\n      \"/etc/crontabs/root\"\n    ]\n  },\n  \"condition\": \"1 of selection_* and not 1 of filter_optional_*\"\n}","reference":["https://github.com/microsoft/MSTIC-Sysmon/blob/f1477c0512b0747c1455283069c21faec758e29d/linux/configs/attack-based/persistence/T1053.003_Cron_Activity.xml","https://pberba.github.io/security/2022/01/30/linux-threat-hunting-for-persistence-systemd-timers-cron/","https://www.elastic.co/security-labs/primer-on-persistence-mechanisms","https://snehbavarva.medium.com/privilege-escalation-techniques-series-linux-cron-jobs-a5b797b424b4"],"false_positives":["Legitimate administrative tasks, package managers, containers, configuration management tools, cloud agents, or system maintenance operations might cause false positives. Apply baselining before deployment."]}}},"__N_SSG":true}