# _next/data/zCRDCjWzEGABGGdCBJ0m0/rules/sigma-6b14bac8-3e3a-4324-8109-42f0546a347f.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2020-10-06T00:00:00.000Z","creation_date":"2020-10-06T00:00:00.000Z","source_siem":"Sigma","source_siem_id":"sigma","source_url":"https://github.com/SigmaHQ/sigma/blob/master/rules/linux/process_creation/proc_creation_lnx_schedule_task_job_cron.yml"},"rule":{"name":"Scheduled Cron Task/Job - Linux","description":"Detects abuse of the cron utility to perform task scheduling for initial or recurring execution of malicious code. Detection will focus on crontab jobs uploaded from the tmp folder.","tags":["Rule Type: Sigma","Data Source: Cross-SIEM","Tactic: Execution","Tactic: Persistence","Technique: T1053.003","OS: Linux","Domain: process creation","SIEM: Sigma"],"severity":"medium","risk_score":47,"type":"sigma","language":"sigma","query":"{\n  \"selection\": {\n    \"Image|endswith\": \"crontab\",\n    \"CommandLine|contains\": \"/tmp/\"\n  },\n  \"condition\": \"selection\"\n}","reference":["https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md"],"false_positives":["Legitimate administration activities"]}}},"__N_SSG":true}