# _next/data/GgAXicA-xLcb6Wga6jhXM/rules/opensearch-6c4e2f43-d94d-4ead-b64d-97e53fa2bd05.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2026-04-17T18:57:30.000Z","source_siem":"OpenSearch","source_siem_id":"opensearch","source_url":"https://github.com/opensearch-project/security-analytics/blob/main/src/main/config/rules/linux/file_event/file_event_lnx_persistence_cron_files.yml"},"rule":{"name":"Persistence Via Cron Files","description":"Detects creation of cron file or files in Cron directories which could indicates potential persistence.","tags":["Rule Type: Sigma","Data Source: OpenSearch","Domain: Cloud","Tactic: Persistence","Technique: T1053.003","OS: Linux","SIEM: OpenSearch"],"severity":"medium","risk_score":47,"type":"sigma","language":"sigma","query":"{\n  \"selection1\": {\n    \"TargetFilename|startswith\": [\n      \"/etc/cron.d/\",\n      \"/etc/cron.daily/\",\n      \"/etc/cron.hourly/\",\n      \"/etc/cron.monthly/\",\n      \"/etc/cron.weekly/\",\n      \"/var/spool/cron/crontabs/\"\n    ]\n  },\n  \"selection2\": {\n    \"TargetFilename|contains\": [\n      \"/etc/cron.allow\",\n      \"/etc/cron.deny\",\n      \"/etc/crontab\"\n    ]\n  },\n  \"condition\": \"1 of selection*\"\n}","reference":["https://github.com/microsoft/MSTIC-Sysmon/blob/f1477c0512b0747c1455283069c21faec758e29d/linux/configs/attack-based/persistence/T1053.003_Cron_Activity.xml"],"false_positives":["Any legitimate cron file."]}}},"__N_SSG":true}