# _next/data/GgAXicA-xLcb6Wga6jhXM/rules/sigma-65d506d3-fcfe-4071-b4b2-bcefe721bbbb.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2023-02-18T00:00:00.000Z","creation_date":"2023-02-18T00:00:00.000Z","source_siem":"Sigma","source_siem_id":"sigma","source_url":"https://github.com/SigmaHQ/sigma/blob/master/rules/macos/process_creation/proc_creation_macos_persistence_via_plistbuddy.yml"},"rule":{"name":"Potential Persistence Via PlistBuddy","description":"Detects potential persistence activity using LaunchAgents or LaunchDaemons via the PlistBuddy utility","tags":["Rule Type: Sigma","Data Source: Cross-SIEM","Tactic: Persistence","Technique: T1543.001","Technique: T1543.004","OS: Macos","Domain: process creation","SIEM: Sigma"],"severity":"high","risk_score":73,"type":"sigma","language":"sigma","query":"{\n  \"selection\": {\n    \"Image|endswith\": \"/PlistBuddy\",\n    \"CommandLine|contains|all\": [\n      \"RunAtLoad\",\n      \"true\"\n    ],\n    \"CommandLine|contains\": [\n      \"LaunchAgents\",\n      \"LaunchDaemons\"\n    ]\n  },\n  \"condition\": \"selection\"\n}","reference":["https://redcanary.com/blog/clipping-silver-sparrows-wings/","https://www.manpagez.com/man/8/PlistBuddy/"],"false_positives":["Unknown"]}}},"__N_SSG":true}