# _next/data/2ZO-91VEXVdEAsqHqXec7/rules/sigma-1bac86ba-41aa-4f62-9d6b-405eac99b485.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2022-02-03T00:00:00.000Z","creation_date":"2022-02-03T00:00:00.000Z","source_siem":"Sigma","source_siem_id":"sigma","source_url":"https://github.com/SigmaHQ/sigma/blob/master/rules/linux/auditd/path/lnx_auditd_systemd_service_creation.yml"},"rule":{"name":"Systemd Service Creation","description":"Detects a creation of systemd services which could be used by adversaries to execute malicious code.","tags":["Rule Type: Sigma","Data Source: Cross-SIEM","Tactic: Persistence","Technique: T1543.002","OS: Linux","SIEM: Sigma"],"severity":"medium","risk_score":47,"type":"sigma","language":"sigma","query":"{\n  \"path\": {\n    \"type\": \"PATH\",\n    \"nametype\": \"CREATE\"\n  },\n  \"name_1\": {\n    \"name|startswith\": [\n      \"/usr/lib/systemd/system/\",\n      \"/etc/systemd/system/\"\n    ]\n  },\n  \"name_2\": {\n    \"name|contains\": \"/.config/systemd/user/\"\n  },\n  \"condition\": \"path and 1 of name_*\"\n}","reference":["https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1543.002/T1543.002.md"],"false_positives":["Admin work like legit service installs."]}}},"__N_SSG":true}