# _next/data/r0OjT_dXJyMb45UcmXflx/rules/opensearch-af202fd3-7bff-4212-a25a-fb34606cfcbe.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2026-05-28T16:55:34.000Z","source_siem":"OpenSearch","source_siem_id":"opensearch","source_url":"https://github.com/opensearch-project/security-analytics/blob/main/src/main/config/rules/linux/builtin/cron/lnx_cron_crontab_file_modification.yml"},"rule":{"name":"Modifying Crontab","description":"Detects suspicious modification of crontab file.","tags":["Rule Type: Sigma","Data Source: OpenSearch","Domain: Cloud","Tactic: Persistence","Technique: T1053.003","OS: Linux","SIEM: OpenSearch"],"severity":"medium","risk_score":47,"type":"sigma","language":"sigma","query":"{\n  \"keywords\": [\n    \"REPLACE\"\n  ],\n  \"condition\": \"keywords\"\n}","reference":["https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.003/T1053.003.md"],"false_positives":["Legitimate modification of crontab"]}}},"__N_SSG":true}