# _next/data/r0OjT_dXJyMb45UcmXflx/rules/sigma-aff229ab-f8cd-447b-b215-084d11e79eb0.json
{"pageProps":{"rule":{"metadata":{"updated_date":"2022-01-01T00:00:00.000Z","creation_date":"2022-01-01T00:00:00.000Z","source_siem":"Sigma","source_siem_id":"sigma","source_url":"https://github.com/SigmaHQ/sigma/blob/master/rules/application/rpc_firewall/rpc_firewall_sasec_lateral_movement.yml"},"rule":{"name":"Remote Schedule Task Lateral Movement via SASec","description":"Detects remote RPC calls to create or execute a scheduled task via SASec","tags":["Rule Type: Sigma","Data Source: Cross-SIEM","Tactic: Execution","Tactic: Persistence","Technique: T1053","Technique: T1053.002","OS: Rpc_firewall","Domain: application","SIEM: Sigma"],"severity":"high","risk_score":73,"type":"sigma","language":"sigma","query":"{\n  \"selection\": {\n    \"EventLog\": \"RPCFW\",\n    \"EventID\": 3,\n    \"InterfaceUuid\": \"378e52b0-c0a9-11cf-822d-00aa0051e40f\",\n    \"OpNum\": [\n      0,\n      1\n    ]\n  },\n  \"condition\": \"selection\"\n}","reference":["https://learn.microsoft.com/en-us/openspecs/windows_protocols/ms-tsch/d1058a28-7e02-4948-8b8d-4a347fa64931","https://github.com/jsecurity101/MSRPC-to-ATTACK/blob/ddd4608fe8684fcf2fcf9b48c5f0b3c28097f8a3/documents/MS-TSCH.md","https://github.com/zeronetworks/rpcfirewall","https://zeronetworks.com/blog/stopping-lateral-movement-via-the-rpc-firewall/"],"false_positives":["Unknown"]}}},"__N_SSG":true}